# systemd/pi-network-admin.service
[Unit]
Description=ESS Datalogger UI Web Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi-network-admin
Group=pi-network-admin
WorkingDirectory=/opt/pi-network-admin
EnvironmentFile=/etc/pi-network-admin/app.env
ExecStart=/opt/pi-network-admin/.venv/bin/gunicorn --bind 0.0.0.0:${PI_ADMIN_PORT} --workers 2 --timeout 60 run:app
ExecStartPost=/usr/bin/bash -c '\
  sleep 1; \
  /usr/bin/nmcli radio wifi on 2>/dev/null; \
  sleep 0.5; \
  /usr/bin/nmcli device set wlan0 managed off 2>/dev/null; \
  sleep 0.5; \
  /usr/bin/nmcli device set wlan0 managed on 2>/dev/null; \
  sleep 0.5; \
  /usr/bin/logger -t pi-network-admin "WiFi interface reset completed after service start"'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
